# build/ports/msmtp/files/msmtp.in
#!/bin/sh
#
# Part of NAS4Free (http://www.nas4free.org).
# Copyright (c) 2012-2018 The NAS4Free Project <info@nas4free.org>.
# All rights reserved.
#
# msmtp configuration
#

# PROVIDE: msmtp
# BEFORE: DAEMON

. /etc/rc.subr
. /etc/configxml.subr

name="msmtp"
desc="msmtp configuration"
load_rc_config "${name}"

# Defaults
msmtp_config=${msmtp_config:-"/var/etc/msmtp.conf"}
msmtprc=${msmtprc:-"/usr/local/etc/msmtprc"}

# Create msmtp.conf configuration
/usr/local/bin/xml sel -t \
	-o "account default" -n \
	-i "string-length(//system/domain) > 0" \
		--if "string-length(//system/hostname) > 0" \
			-v "concat('domain ',//system/hostname,'.',//system/domain)" -n \
		--else \
			-v "concat('domain ',//system/domain)" -n \
		--break \
	-b \
	-i "string-length(//system/email/server) > 0" \
		-v "concat('host ',//system/email/server)" -n \
	-b \
	-i "string-length(//system/email/port) > 0" \
		-v "concat('port ',//system/email/port)" -n \
	-b \
	-o "protocol smtp" -n \
	-i "string-length(//system/email/from) > 0" \
		-v "concat('from ',//system/email/from)" -n \
	-b \
	-i "count(//system/email/auth) > 0" \
		-i "string-length(//system/email/authmethod) > 0" \
			-v "concat('auth ',//system/email/authmethod)" -n \
		-b \
		-i "string-length(//system/email/username) > 0" \
			-v "concat('user ',//system/email/username)" -n \
		-b \
		-i "string-length(//system/email/password) > 0" \
			-v "concat('password ',//system/email/password)" -n \
		-b \
	-b \
	--if "count(//system/email/starttls) > 0" \
		-o "tls on" -n \
		-o "tls_starttls on" -n \
		--if "string-length(//system/email/tls_certcheck) > 0" \
			-v "concat('',//system/email/tls_certcheck)" -n \
			--if "contains(//system/email/tls_certcheck,'tls_certcheck on')" \
				--if "string-length(//system/email/tls_trust_file) > 0" \
					-v "concat('tls_trust_file ',//system/email/tls_trust_file)" -n \
				--else \
					-o "tls_trust_file /usr/local/etc/ssl/cert.pem" -n \
				--break \
			--break \
		--else \
			-o "tls_certcheck off" -n \
		--break \
	--elif "//system/email/security[. != 'none']" \
		-o "tls on" -n \
		-o "tls_starttls off" -n \
		--if "string-length(//system/email/tls_certcheck) > 0" \
			-v "concat('',//system/email/tls_certcheck)" -n \
			--if "contains(//system/email/tls_certcheck,'tls_certcheck on')" \
				--if "string-length(//system/email/tls_trust_file) > 0" \
					-v "concat('tls_trust_file ',//system/email/tls_trust_file)" -n \
				--else \
					-o "tls_trust_file /usr/local/etc/ssl/cert.pem" -n \
				--break \
			--break \
		--else \
			-o "tls_certcheck off" -n \
		--break \
	--else \
		-o "tls off" -n \
		-o "tls_starttls off" -n \
	--break \
	-o "syslog LOG_MAIL" -n \
	${configxml_file} | /usr/local/bin/xml unesc > ${msmtp_config}

# Set user permissions
chmod 0600 ${msmtp_config}
# for PHP mail()
cp ${msmtp_config} ${msmtprc}
chown root:www ${msmtprc}
chmod 0640 ${msmtprc}
